# etc/profile.d/ublue-just.sh
# set alias for global justfile usage
alias ublue="just --justfile /etc/justfile"